# JavaScript/NodeCode/AwsBasicS3/AwsBasic.conf
# This is an upstart script: http://upstart.ubuntu.com/index.html
description "a script to keep node.js server in memory even after rebooting"
author      "Charle Calvert - http://www.elvenware.com/charlie"


# Start after all drives mounted
start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    export HOME="/root"
    export JSOBJECTS=/home/charlie/Git/JsObjects
    export PYTHONPATH=${PYTHONPATH}:$JSOBJECTS/Python/:$JSOBJECTS/Python/Utils/:$JSOBJECTS/Python/RegEx/:

# The following assumes nodejs is in /usr/bin
    exec /usr/bin/nodejs /home/charlie/bin/AwsBasic/app.js >> /var/log/node.log 2>&1
# exec /home/charlie/bin/awsbasic >> /var/log/node.log 2>&1
end script

post-start script
   # Optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "node.js has started!"
end script

